# open_source_ami/flasher-1.19.0-ARM/package/data/flasher
# /etc/init.d/flasher: start flasher
#
# chkconfig: 2345 10 90
#Runlevel : 3 = S80
PATH=/bin:/usr/bin:/sbin:/usr/sbin

PATHFUL_FLASHER="/usr/local/bin/flasher"



test -f /usr/local/bin/flasher || exit 0


# Options for start/restart the daemons
#


#


case "$1" in
  start)
    echo "Starting Firmware Upgrade Utility: flasher."
    start-stop-daemon --start --quiet --exec $PATHFUL_FLASHER 
    ;;
  stop)
    echo "Stopping Firmware Upgrade Utility: flasher."
    start-stop-daemon --stop --quiet --exec $PATHFUL_FLASHER
    ;;
   *)
    echo "Usage: /etc/init.d/flasher {start|stop}"
    exit 1
esac

exit 0
