# open_source_ami/sys_base-2.41.0-ARM/package/data/etc/network/if-pre-up.d/01-pre-up-script
#!/bin/sh
echo "Interface $IFACE is going to go up"
exit 0
