# source/sys_base_ex-2.56.0-src/data/defaults/etc/network/if-up.d/01-up-script
#!/bin/sh
echo "Interface $IFACE is up"
exit 0
