# source/sys_base_ex-2.19.0-src/data/defaults/etc/network/if-pre-up.d/01-pre-up-script
#!/bin/sh
echo "Interface $IFACE is going to go up"
echo 0 > /proc/sys/net/ipv6/conf/default/accept_ra  # Quanta
echo 0 > /proc/sys/net/ipv6/conf/default/autoconf   # Quanta
exit 0
